# virl/std/files/virl-webmux.conf
start on runlevel [235] or starting virl-uwm
stop on runlevel [06] or stopping virl-uwm

setuid virl
setgid virl

exec /usr/local/bin/virl_webmuxd
#pre-stop exec initctl emit kill-webmux
